# deploy/systemd/paramant-relay-finance.service
[Unit]
Description=PARAMANT Relay - finance
After=network.target

[Service]
Environment=PORT=3003
Environment=SECTOR=finance
Environment=ADMIN_TOKEN=REPLACE_ADMIN_TOKEN
Environment=USERS_FILE=/home/paramant/relay-finance/users.json
Environment=TOTP_SECRET=REPLACE_TOTP_SECRET
Environment=RESEND_API_KEY=REPLACE_RESEND_KEY
Environment=LOG_LEVEL=info
Type=simple
User=paramant
WorkingDirectory=/home/paramant/relay-finance
ExecStart=/usr/bin/node relay.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
