# deploy/systemd/paramant-relay-health.service
[Unit]
Description=PARAMANT Relay - health
After=network.target

[Service]
Environment=TOTP_SECRET=REPLACE_TOTP_SECRET
Environment=RESEND_API_KEY=REPLACE_RESEND_KEY
Type=simple
User=paramant
WorkingDirectory=/home/paramant/relay-health
ExecStart=/usr/bin/node relay.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
