# deploy/systemd/paramant-relay-health.service
[Unit]
Description=PARAMANT Relay - health
After=network.target

[Service]
Environment=TOTP_SECRET=GGM3VTGW4DVPN7QHABQ7JXZP6RLUFSO7
Environment=RESEND_API_KEY=re_xxxx_redacted
Type=simple
User=paramant
WorkingDirectory=/home/paramant/relay-health
ExecStart=/usr/bin/node relay.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
